# scripts/linux/hostlink.service
[Unit]
Description=hostlink
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/usr/bin/
EnvironmentFile=/etc/hostlink/hostlink.env
ExecStart=/usr/bin/hostlink
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target

